# docker/s6-rc.d/grafana-data-init/up
#!/command/execlineb -P
with-contenv
/scripts/init-grafana-data.sh
