# docker/s6-rc.d/collector/run
#!/command/execlineb -P
cd /app
/usr/local/bin/bun run src/index.ts
