# docker/s6-rc.d/postgres-data-init/up
#!/command/execlineb -P
with-contenv
/scripts/init-postgres-data.sh
