# docker/s6-rc.d/init-env/up
#!/command/execlineb -P
with-contenv
/scripts/validate-env.sh
