# docker/s6-rc.d/postgres-init/up
#!/command/execlineb -P
with-contenv
/scripts/init-postgres.sh
