# docker/s6-rc.d/postgres/run
#!/command/execlineb -P
s6-setuidgid postgres
postgres -D /var/lib/postgresql/data
